# enigma/enigma
# Copyright (c) 2016 NIST
# All Rights Reserved.
# Author: Junkai Huang

#!/bin/sh

#
### BEGIN INIT INFO
# Provides:          Enigma
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:      $network $named $time
# Should-Stop:       $network $named $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the SpringX  server daemon
### END INIT INFO
#


#set -e

FILENAME=pid/process.pid

. /lib/lsb/init-functions


springx_stop() {
    log_daemon_msg "Stopping the Enigma-services"
    if [ ! -s $FILENAME ]; then
        log_failure_msg "Enigma-service is not active."
        log_end_msg 1
        return 0
    fi

    while read line
    do
         kill -9  $line
         sed -i '/'$line'/d' pid/process.pid
    done< $FILENAME
    sleep 2
    log_end_msg 0
}

springx_start() {
    log_daemon_msg "Starting the Enigma-services"
    if [ -s $FILENAME ]; then 
        log_failure_msg "Enigma-service is running."
        log_end_msg 1
        return 0
    fi
    #python3 enigma.py
    python3 service.py&
    log_end_msg 0
}

springx_status() {
    if [ -s $FILENAME ]; then 
        log_warning_msg "Enigma-service is running."
    else
        log_warning_msg "Enigma-service is not active." 
   fi 
}

case $1 in
    restart)
        springx_stop;springx_start;;
    start)
        springx_start;;
    status)
        springx_status;;
    stop)
        springx_stop;;
    *) 
       log_failure_msg "$1: Error: The action is not supported.";;
esac

exit 0
